# dist/supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
user=tpot
autorestart=true

[program:kibana]
command=/usr/share/kibana/bin/kibana
autorestart=true

[program:logstash]
command=/bin/bash -c 'cp /data/elk/logstash/conf/*.conf /etc/logstash/conf.d/ || true && exec /usr/share/logstash/bin/logstash -f /etc/logstash/conf.d/logstash.conf'
autorestart=true

[program:head]
directory=/opt/elasticsearch-head
command=/usr/local/bin/grunt server
autorestart=true

